# shellscript.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/home/docker/script.sh

[Install]
WantedBy=multi-user.target

